# services/linux/claude-limit-monitor.service
[Unit]
Description=Claude Unified Limit Monitor
After=network.target

[Service]
Type=simple
User=tallest
WorkingDirectory=/home/tallest/Claude
Environment="CLAUDE_HOME=/home/tallest/Claude"
ExecStart=/usr/bin/python3 /home/tallest/Claude/scripts/claude-limit-daemon.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-limit-monitor

[Install]
WantedBy=multi-user.target
